# lib/systemd/user/yt-dlp-update.service
[Unit]
Description=Update %S/yt-dlp.git

[Service]
ExecStart=/usr/local/libexec/update-yt-dlp %S/yt-dlp.git
Restart=on-failure
Type=oneshot
